# etc/crontabs/root
# do daily/weekly/monthly maintenance
# min	hour	day	month	weekday	command
*/5	*	*	*	*	run-parts /etc/periodic/ddns
